# monitoring/03_sentinelle_unifiee_multijails_corrigee/systemd/sentinelle.service
[Unit]
Description=Sentinelle - Protection anti brute-force locale (equivalent Fail2Ban en Bash)
After=network.target rsyslog.service

[Service]
Type=simple
ExecStart=/opt/sentinelle/sentinelle.sh
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
